# scripts/systemd/sejfa-loop-runner.service
[Unit]
Description=SEJFA Loop Runner (polls queue, triggers Ralph Loop)
After=network.target

[Service]
Type=simple
# Loop runner runs on Mac, not Hetzner — this unit is for reference.
# Update User/WorkingDirectory/paths for your Mac or deployment target.
User=sejfa
WorkingDirectory=/opt/sejfa
Environment=LOOP_RUNNER_BACKEND_URL=https://sejfa.example.com
Environment=LOOP_RUNNER_REPO_DIR=/opt/sejfa
Environment=LOOP_RUNNER_POLL_INTERVAL=10
Environment=PYTHONPATH=/opt/sejfa/services/loop-engine/src
ExecStart=/opt/sejfa/venv/bin/python -m loop_engine.runner
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sejfa-loop

[Install]
WantedBy=multi-user.target
