# scripts/systemd/sejfa-monitor-api.service
[Unit]
Description=SEJFA Monitor API (FastAPI on :8100)
After=network.target

[Service]
Type=simple
User=sejfa
WorkingDirectory=/opt/sejfa
EnvironmentFile=/opt/sejfa/.env
ExecStart=/opt/sejfa/venv/bin/uvicorn monitor.api:app --host 0.0.0.0 --port 8100
Environment=PYTHONPATH=/opt/sejfa/services/monitor-api/src
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sejfa-monitor

[Install]
WantedBy=multi-user.target
